# deploy/systemd/hub-backup.timer
[Unit]
Description=Hub nightly backup timer
Requires=hub-backup.service

[Timer]
# 03:30 local time every day. Persistent=true catches missed runs after reboot.
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
